# contrib/com.ebay-gpu-hunter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- Unique label — used by launchctl to identify this job -->
  <key>Label</key>
  <string>com.ebay-gpu-hunter</string>

  <!-- Command to run. Adjust paths if your project or venv lives elsewhere. -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/glisk/Workdir/ebay-hunter/.venv/bin/python</string>
    <string>/Users/glisk/Workdir/ebay-hunter/gpu_hunt.py</string>
    <string>--report</string>
  </array>

  <!-- Working directory — gpu-report.md and cache/ are written relative to this -->
  <key>WorkingDirectory</key>
  <string>/Users/glisk/Workdir/ebay-hunter</string>

  <!-- Run every 2 hours (7200 seconds) -->
  <key>StartInterval</key>
  <integer>7200</integer>

  <!-- Stdout and stderr go to separate log files from workstation hunter -->
  <key>StandardOutPath</key>
  <string>/Users/glisk/Workdir/ebay-hunter/cache/gpu-launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/glisk/Workdir/ebay-hunter/cache/gpu-launchd-error.log</string>

  <!-- Restart automatically if it crashes -->
  <key>KeepAlive</key>
  <false/>

  <!-- Run even if the job was missed while the machine was asleep -->
  <key>RunAtLoad</key>
  <false/>

</dict>
</plist>
